# systemd/user/alacritty.service
[Unit]
Description=Start Alacritty

[Service]
Type=simple
ExecStart=/run/current-system/sw/bin/alacritty

[Install]
WantedBy=default.target
