# systemd/user/vscode.service
[Unit]
Description=Start VSCode

[Service]
Type=simple
ExecStart=/run/current-system/sw/bin/code

[Install]
WantedBy=default.target
