# systemd/user/waybar.service
[Unit]
Description=Start Waybar

[Service]
Type=simple
ExecStart=/run/current-system/sw/bin/waybar

[Install]
WantedBy=default.target
